# deploy/maker-dashboard-update.service
[Unit]
Description=Pull, verify, and deploy the latest Maker Dashboard image
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/maker-dashboard-update.sh
SyslogIdentifier=maker-dashboard-update
